# systemd/kiosk-display.service
[Unit]
Description=KioskPi Display (Cage + Chromium)
After=network.target systemd-user-sessions.service getty@tty1.service
Conflicts=getty@tty1.service

[Service]
Type=simple
User=kiosk
PAMName=login
Environment=WLR_LIBINPUT_NO_DEVICES=1
Environment=XDG_RUNTIME_DIR=/run/user/%U
Environment=WAYLAND_DISPLAY=wayland-1
StandardOutput=journal
StandardError=journal
ExecStart=/usr/bin/cage -s -- /usr/bin/chromium-browser \
    --enable-features=VaapiVideoDecoder,VaapiVideoEncoder \
    --enable-gpu-rasterization \
    --enable-zero-copy \
    --disable-software-rasterizer \
    --ignore-gpu-blocklist \
    --disable-dev-shm-usage \
    --disable-extensions \
    --disable-component-update \
    --disable-background-networking \
    --disable-sync \
    --no-first-run \
    --kiosk \
    --noerrdialogs \
    --disable-infobars \
    --remote-debugging-port=9222 \
    --user-data-dir=/home/kiosk/.config/chromium-kiosk \
    about:blank
Restart=on-failure
RestartSec=5
StandardInput=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes

[Install]
WantedBy=multi-user.target
